# docker/carla_ros_bridge/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
loglevel=error

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:carla_bridges]
command=/bin/bash -c "ros2 launch carla_ros_bridge carla_ros_bridge_with_example_ego_vehicle.launch.py host:=CONTAINER_NAME_carla timeout:=200"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true


[program:code-server]
command=/home/docker/colcon_ws/start_code_server.sh
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/tmp/code-server.log
redirect_stderr=true
